# packaging/systemd/pixelensd.service
[Unit]
Description=Pixelens daemon for visual search and OCR
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/pixelensd
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
